# Buraq/Centos/crontab
@reboot sleep 20;/usr/bin/vmrun start "/vm/Odoo/VM-Name.vmx" nogui 
