# Buraq/VM/postgres/crontab
45 */4 * * * /var/lib/postgresql/backup.sh > /dev/null 2>&1
